# custom-components/services/fix-ssh-sudoers.service
[Unit]
Description=Ghettoblaster SSH and Sudoers Fix (PERMANENT)
After=multi-user.target
# REMOVED: After=moode-startup.service (can cause boot hang - moode-startup may not exist)
Wants=multi-user.target
DefaultDependencies=no

[Service]
Type=oneshot
# CRITICAL: Enable and start SSH in multiple ways
# All commands have timeouts to prevent hangs
ExecStart=/bin/bash -c 'timeout 3 bash -c "echo \"andre ALL=(ALL) NOPASSWD: ALL\" > /etc/sudoers.d/andre && chmod 0440 /etc/sudoers.d/andre" 2>/dev/null || true; timeout 2 visudo -c -f /etc/sudoers.d/andre 2>/dev/null || true; timeout 2 systemctl enable ssh 2>/dev/null || timeout 2 systemctl enable sshd 2>/dev/null || true; timeout 2 systemctl start ssh 2>/dev/null || timeout 2 systemctl start sshd 2>/dev/null || true; timeout 2 touch /boot/firmware/ssh 2>/dev/null || timeout 2 touch /boot/ssh 2>/dev/null || true; exit 0'
RemainAfterExit=yes
TimeoutStartSec=5
TimeoutStopSec=2

[Install]
WantedBy=multi-user.target

