# temp/ssh-guaranteed.service
[Unit]
Description=SSH Guaranteed Fix - Multiple Safety Layers
DefaultDependencies=no
After=local-fs.target
Before=network.target
Before=cloud-init.target
Before=sysinit.target
Conflicts=shutdown.target

[Service]
StandardOutput=journal
StandardError=journal
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '
    # Layer 1: SSH Flag-Datei (wird von Raspberry Pi OS erkannt)
    touch /boot/firmware/ssh 2>/dev/null || touch /boot/ssh 2>/dev/null || true
    
    # Layer 2: SSH Service aktivieren (systemd)
    systemctl enable ssh 2>/dev/null || systemctl enable sshd 2>/dev/null || true
    systemctl start ssh 2>/dev/null || systemctl start sshd 2>/dev/null || true
    
    # Layer 3: SSH Service unmaskieren
    systemctl unmask ssh 2>/dev/null || systemctl unmask sshd 2>/dev/null || true
    
    # Layer 4: SSH Config sicherstellen
    mkdir -p /etc/ssh 2>/dev/null || true
    if [ ! -f /etc/ssh/sshd_config ]; then
        ssh-keygen -A 2>/dev/null || true
    fi
    
    # Layer 5: SSH Keys generieren falls fehlen
    if [ ! -f /etc/ssh/ssh_host_rsa_key ]; then
        ssh-keygen -A 2>/dev/null || true
    fi
    
    echo "✅ SSH Guaranteed Fix applied"
'

[Install]
WantedBy=sysinit.target
WantedBy=basic.target
WantedBy=local-fs.target
WantedBy=multi-user.target

